# launchd/com.proton.pass-cli.ssh-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.proton.pass-cli.ssh-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mrossi/dotfiles/launchd/proton-pass-ssh-agent.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/mrossi/Library/Logs/proton-pass-ssh-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mrossi/Library/Logs/proton-pass-ssh-agent.log</string>
</dict>
</plist>
